# haproxy.cron
#!/bin/bash

#
# Renew certificates from LetsEncrypt
#

if [ ! -z "$HAPROXY_LETSENCRYPT" ]; then
	/usr/local/sbin/certbot-renew >/dev/null
fi


#
# Reload any configuration changes
#

if [ ! -z "$HAPROXY_DAILY_RELOAD" ]; then
        systemctl reload haproxy
fi


#
# Wait for reload to finish then update OCSP stapling
#

if [ ! -z "$HAPROXY_LETSENCRYPT" ] && [ ! -z "$HAPROXY_LETSENCRYPT_OCSP" ]; then
	sleep 15
	/usr/local/sbin/certbot-ocsp >/dev/null
fi
